# docs/matrix-dirwatch.service
[Unit]
Description=Watch directories and post changes to matrix

[Service]
Type=simple
ExecStart=matrix-dirwatch.sh
RestartSec=5s
